# home/private_Library/private_LaunchAgents/com.germanfrelo.chezmoi-drift.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.germanfrelo.chezmoi-drift</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-c</string>
		<string>status=$(timeout 10 /opt/homebrew/bin/chezmoi status 2>/dev/null); [[ -n "$status" ]] &amp;&amp; /usr/bin/osascript -e 'display notification "Run chezmoi diff --use-builtin-diff to review." with title "⚠ chezmoi drift detected"'</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>StandardOutPath</key>
	<string>/dev/null</string>
	<key>StandardErrorPath</key>
	<string>/dev/null</string>
</dict>
</plist>
